# server/crontab/ybCrontab

*/1 * * * * php /usr/local/nginx/html/yb-job/ybTask demo
